# dist/mravatar.service
[Unit]
Description=Mravatar
After=network.target

[Service]
User=mravatar
Group=mravatar
WorkingDirectory=/home/mravatar/mravatar
Environment="PATH=/home/mravatar/.local/bin"
ExecStart=/home/mravatar/.local/bin/gunicorn mravatar:app --bind 0.0.0.0:8000

[Install]
WantedBy=multi-user.target